# resources/python/linux/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/bachtiarpanjaitan/Documents/bachtiarpanjaitan/ic-php/resources/python/linux/bin/conda";
setenv _CONDA_EXE "/home/bachtiarpanjaitan/Documents/bachtiarpanjaitan/ic-php/resources/python/linux/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/home/bachtiarpanjaitan/Documents/bachtiarpanjaitan/ic-php/resources/python/linux/bin/python";
setenv _CONDA_ROOT "/home/bachtiarpanjaitan/Documents/bachtiarpanjaitan/ic-php/resources/python/linux";

source "/home/bachtiarpanjaitan/Documents/bachtiarpanjaitan/ic-php/resources/python/linux/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"